# airootfs/etc/profile.d/pallas.sh
#!/bin/bash
# PallasOS login banner

if [[ $- == *i* ]]; then
    echo ''
    echo '                                         .-+############+-.                                         '
    echo '                                     .########################.                                     '
    echo '                                  -#####+-.              .-+#####-                                  '
    echo '                      +########-#####-                        .#####-########+                      '
    echo '                    .###-..-+#####+                              -#####+-..-###-                    '
    echo '                   .##-     -###-                                  -###-     -##.                   '
    echo '                   +## +##+###+                                      +###+##+ ##+                   '
    echo '                   +#+ ######.                                        .######.+#+                   '
    echo '                   +## #####                                            #####.##+                   '
    echo '                   .##.-###.                                             ###+.##-                   '
    echo '                   .######-                                              -######                    '
    echo '                    -#####                                                +####-                    '
    echo '                     ####.                 .-.        .--                 .####                     '
    echo '                      ###                   +##+    +###                   ###                      '
    echo '                     .##+     +############-  +##+-###. -############+     +##.                     '
    echo '                     -##-   +###   .#########+  +###. -#########.   +###   -##-                     '
    echo '                     +##- -##+.      -+++-  ##+   .  +##. -+++-       +##- -##+                     '
    echo '                     -##-+##.               +#+      +#+               .##+-##+                     '
    echo '                     -#####.     +########  +#+      +#+  #########      #####-                     '
    echo '                      ####.    +##+         +#+      +#+         -##+    .####.                     '
    echo '                      +###    +##.          +##########+           ###.   +##+                      '
    echo '                      .###-  ##+             -++####++-             +##  -###.                      '
    echo '                       -###-##+                .####.                +##-###-                       '
    echo '                        -#####.                 -##-                 .#####+                        '
    echo '                          .###-    #+       -+########+-       +#    -###.                          '
    echo '                            -###+  ############-    -############  +###+                            '
    echo '                              -####+#######.            .#######+####-                              '
    echo '                                 +####- +#####++-..--+#####+.-####+                                 '
    echo '                                   .-.     -+##########+-     .+.                                   '
    echo ''
    echo '                                           PallasOS Alpha'
    echo '                                         Based on Arch Linux'
    echo ''

    # Set PallasOS bash prompt
    if [ "$(id -u)" -eq 0 ]; then
        PS1='\[\e[1;31m\]\u@pallas\[\e[0m\]:\[\e[1;34m\]\w\[\e[0m\]# '
    else
        PS1='\[\e[1;32m\]\u@pallas\[\e[0m\]:\[\e[1;34m\]\w\[\e[0m\]\$ '
    fi
fi

